# recipes-core/runit-base-services/files/socklogd/sv/syslog/run
#!/bin/sh
exec 2>&1
exec chpst -Unobody:nogroup socklog unix /dev/log
