# recipes-core/eudev/files/sv/udevd/run
#!/bin/sh
udevadm control --exit
exec udevd 2>&1
