# recipes-core/dbus/files/sv/messagebus/run
#!/bin/sh
[ ! -d /run/dbus ] && install -m 0755 -g 996 -o 997 -d /run/dbus
dbus-uuidgen --ensure
exec dbus-daemon --system --nofork --nopidfile

