# external/openembedded-layer/recipes-connectivity/dnsmasq/files/sv/dnsmasq/run
#!/bin/sh
exec dnsmasq -d -C /etc/dnsmasq.conf 2>&1
