# external/networking-layer/recipes-support/ntp/files/sv/isc-ntpd/run
#!/bin/sh

# Pick up our core-services/run files helpers...we MIGHT need one of them.
. /etc/runit/functions

# Check for the presence of an "ntpd" user.  If one isn't present in the system,
# warn the user and stop the service.  That way the supervisory server won't
# keep trying to restart this one.
check_for_user "ntpd"
if [ $? -eq 1 ] ; then
    # Tell our user this isn't working and exit the service
    svc_log "ntpd user not present on system.  Stopping service."
    svc_stop
fi

# Run the service...
exec ntpd -g -u ntpd:ntpd -n >/dev/null 2>&1
