# recipes-core/runit-base-services/files/sv/klog/run
#!/bin/sh

. /etc/runit/syslog-helper

exec klogd -a /sbin/klogd -- -n 2>&1
