# external/networking-layer/recipes-daemons/ptpd/files/sv/ptpd/run
#!/bin/sh

# We want to wait to 'syslog' being "up" before starting this- so we can get
# feedback on a fail of config, etc.
. /etc/runit/functions
svc_running syslog

# Pick up our core-services/run files helpers...we MIGHT need one of them.
. /etc/runit/functions

# Check to see if we have config changes...
[ -e /etc/default/ptpd ] && . /etc/default/ptpd

# Double check to see we have sane settings and if not, default them
# in turn...
[ -z "$PTPD_PID_FILE" ] && PTPD_PID_FILE="/var/run/ptpd.pid"
[ -z "$PTPD_STATUS_FILE" ] && PTPD_STATUS_FILE="/var/run/ptpd.status"
[ -z "$PTPD_CONFIG_FILE" ] && PTPD_CONFIG_FILE="/etc/ptpd2.conf"

# Check to see if we have the config file- no config, no start.  Tell the world and stop
# ourselves...
if [ ! -e $PTPD_CONFIG_FILE ] ; then
    svc_log "ptp config ($PTPD_CONFIG_FILE) is not present.  Stopping service"
    svc_stop
fi

# Launch the daemon...we've got everything.
exec ptpd2 --global:lock_file=$PTPD_PID_FILE --global:status_file=$PTPD_STATUS_FILE -c $PTPD_CONFIG_FILE 2>&1
