# recipes-navigation/gpsd/files/sv/gpsd/run
#!/bin/sh
[ -r /etc/default/gpsd ] && . /etc/default/gpsd
exec /usr/sbin/gpsd -N -F /run/gpsd.sock $OPTS ${DEV:=/dev/gps0}

