# external/openembedded-layer/recipes-navigation/gpsd/files/sv/gpsd/run
#!/bin/sh
[ -r /etc/default/gpsd ] && . /etc/default/gpsd
[ -z "SOCKET" ] && SOCKET="/run/gpsd.sock"
[ -z "DEVICES" ] && DEVICES="/dev/gps0"
exec /usr/sbin/gpsd -N -F $SOCKET $OPTIONS ${DEVICES}

