# recipes-core/dbus/files/sv/dbus/run
# Keeping it simple... We don't need to wait for anything...
exec 2>&1 /usr/bin/dbus-daemon
