# recipes-extended/rpcbind/files/sv/rpcbind/run
#!/bin/sh

OPTIONS=""
if [ -f /etc/default/rpcbind ]
then
    . /etc/default/rpcbind
elif [ -f /etc/rpcbind.conf ]
then
    . /etc/rpcbind.conf
fi

exec /usr/sbin/rpcbind -f $OPTIONS 2>&1
