# recipes-connectivity/wpa-supplicant/files/sv/wpa_supplicant/run
#!/bin/sh
exec wpa_supplicant -Dnl80211 -iwlan0 -c/etc/wpa_supplicant.conf

