# recipes-connectivity/hostapd/files/sv/hostapd/run
#!/bin/sh
exec hostapd /etc/hostapd/hostapd.conf 2>&1
