# recipes-core/runit-base-services/files/sv/syslog/run
#!/bin/sh

. /etc/runit/syslog-helper

exec syslogd -a /sbin/syslogd -- -n $SYSLOG_ARGS 2>&1

