# recipes-connectivity/bluez5/files/sv/bluetooth/run
#!/bin/sh
sv check messagebus > /dev/null || exit 1
exec /usr/libexec/bluetooth/bluetoothd -n > /dev/null

