# recipes-connectivity/bluez5/files/sv/bluetooth/run
# Keeping it simple...
exec 2>&1 /usr/bin/bluetooth/bluetoothd
