# recipes-connectivity/bluez5/files/sv/bluetooth/run
#!/bin/sh
sv check messagebus || exit 1
exec /usr/libexec/bluetooth/bluetoothd -n 2>&1

